# music-player.service
[Unit]
Description=Raspberry Pi Music Player Daemon
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=skayflakes
Group=skayflakes
WorkingDirectory=/home/skayflakes/music-player
ExecStart=/usr/bin/python3 /home/skayflakes/music-player/player.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

