# rodrigo_radio.service
[Unit]
Description=Raspberry Pi Music Player Daemon
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=skayflakes
Group=skayflakes
WorkingDirectory=/home/skayflakes/rodrigo_radio
ExecStart=/usr/bin/python3 /home/skayflakes/rodrigo_radio/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rodrigo_radio

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

